# deploy/groww-trader.service
[Unit]
Description=Optimized Groww Trader Service
After=network.target

[Service]
User=growwtrader
Group=growwtrader
WorkingDirectory=/opt/groww-trader
Environment="PATH=/opt/groww-trader/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/groww-trader/venv/bin/python /opt/groww-trader/optimized_groww_trader.py
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Logging
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=groww-trader

[Install]
WantedBy=multi-user.target
